# roles/duplicity/templates/backup-job.service.j2
[Unit]
Description=Duplicity backup of home directory

[Service]
ExecStart={{ ansible_user_dir }}/bin/backup

[Install]
WantedBy=default.target
